# bin/robot-armhf-install/robot.service
[Unit]
Description=Robot service
After=network.target
Wants=network.target

[Service]
ExecStart=/opt/robot

[Install]
WantedBy=network.target